# docker/nest/cron.d/schedule
*/15 * * * * cd /var/www/nest || npm run parse check-parser
